# deploy/systemd/global-solutions-video-worker.service
[Unit]
Description=YeshuaLife Global Solutions video processing worker
After=network.target

[Service]
Type=simple
WorkingDirectory=/var/www/yeshualife
ExecStart=/root/n/bin/python manage.py process_global_solutions_videos_worker --interval-seconds 20 --limit 1
Restart=always
RestartSec=5

# Hardening (safe defaults; adjust if needed)
PrivateTmp=true
NoNewPrivileges=true

[Install]
WantedBy=multi-user.target

